# conf/dosbox.conf
[program:dosbox]
command=dosbox -fullscreen 
autorestart=true
